# test/ullget.conf
[supervisord]
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor
minfds=100000
minprocs=100000

[group:webservers]
programs=backsoc1,backsoc2,backsoc3,backsoc4,backsoc5,backsoc6,backsoc7,backsoc8,websoc1,websoc2,websoc3,websoc4,websoc5,websoc6,websoc7,websoc8

[program:backsoc1]
command=python3 server.py
directory=/home/ullget/dist/ullget/newsocket/server_back/
user=root
stdout_logfile=/home/ullget/dist/ullget/newsocket/server_back/logs/server.log
numprocs=1
numprocs_start=4
startretries=3
autorestart=true
redirect_stderr=true
loglevel=debug

[program:backsoc2]
command=python3 server.py
directory=/home/ullget/dist/ullget/newsocket/server_back/
user=root
stdout_logfile=/home/ullget/dist/ullget/newsocket/server_back/logs/server.log
numprocs=1
numprocs_start=4
startretries=3
autorestart=true
redirect_stderr=true
loglevel=debug

[program:backsoc3]
command=python3 server.py
directory=/home/ullget/dist/ullget/newsocket/server_back/
user=root
stdout_logfile=/home/ullget/dist/ullget/newsocket/server_back/logs/server.log
numprocs=1
numprocs_start=4
startretries=3
autorestart=true
redirect_stderr=true
loglevel=debug


[program:backsoc4]
command=python3 server.py
directory=/home/ullget/dist/ullget/newsocket/server_back/
user=root
stdout_logfile=/home/ullget/dist/ullget/newsocket/server_back/logs/server.log
numprocs=1
numprocs_start=4
startretries=3
autorestart=true
redirect_stderr=true
loglevel=debug

[program:backsoc5]
command=python3 server.py
directory=/home/ullget/dist/ullget/newsocket/server_back/
user=root
stdout_logfile=/home/ullget/dist/ullget/newsocket/server_back/logs/server.log
numprocs=1
numprocs_start=4
startretries=3
autorestart=true
redirect_stderr=true
loglevel=debug

[program:backsoc6]
command=python3 server.py
directory=/home/ullget/dist/ullget/newsocket/server_back/
user=root
stdout_logfile=/home/ullget/dist/ullget/newsocket/server_back/logs/server.log
numprocs=1
numprocs_start=4
startretries=3
autorestart=true
redirect_stderr=true
loglevel=debug

[program:backsoc7]
command=python3 server.py
directory=/home/ullget/dist/ullget/newsocket/server_back/
user=root
stdout_logfile=/home/ullget/dist/ullget/newsocket/server_back/logs/server.log
numprocs=1
numprocs_start=4
startretries=3
autorestart=true
redirect_stderr=true
loglevel=debug

[program:backsoc8]
command=python3 server.py
directory=/home/ullget/dist/ullget/newsocket/server_back/
user=root
stdout_logfile=/home/ullget/dist/ullget/newsocket/server_back/logs/server.log
numprocs=1
numprocs_start=4
startretries=3
autorestart=true
redirect_stderr=true
loglevel=debug


[program:websoc1]
command=python3 server.py -p 9001
directory=/home/ullget/dist/ullget/newsocket/server_socket/
user=root
stdout_logfile=/home/ullget/dist/ullget/newsocket/server_socket/logs/server.log
numprocs=1
numprocs_start=4
startretries=3
autorestart=true
redirect_stderr=true
loglevel=debug


[program:websoc2]
command=python3 server.py -p 9002
directory=/home/ullget/dist/ullget/newsocket/server_socket/
user=root
stdout_logfile=/home/ullget/dist/ullget/newsocket/server_socket/logs/server.log
numprocs=1
numprocs_start=4
startretries=3
autorestart=true
redirect_stderr=true
loglevel=debug

[program:websoc3]
command=python3 server.py -p 9003
directory=/home/ullget/dist/ullget/newsocket/server_socket/
user=root
stdout_logfile=/home/ullget/dist/ullget/newsocket/server_socket/logs/server.log
numprocs=1
numprocs_start=4
startretries=3
autorestart=true
redirect_stderr=true
loglevel=debug

[program:websoc4]
command=python3 server.py -p 9004
directory=/home/ullget/dist/ullget/newsocket/server_socket/
user=root
stdout_logfile=/home/ullget/dist/ullget/newsocket/server_socket/logs/server.log
numprocs=1
numprocs_start=4
startretries=3
autorestart=true
redirect_stderr=true
loglevel=debug

[program:websoc5]
command=python3 server.py -p 9005
directory=/home/ullget/dist/ullget/newsocket/server_socket/
user=root
stdout_logfile=/home/ullget/dist/ullget/newsocket/server_socket/logs/server.log
numprocs=1
numprocs_start=4
startretries=3
autorestart=true
redirect_stderr=true
loglevel=debug

[program:websoc6]
command=python3 server.py -p 9006
directory=/home/ullget/dist/ullget/newsocket/server_socket/
user=root
stdout_logfile=/home/ullget/dist/ullget/newsocket/server_socket/logs/server.log
numprocs=1
numprocs_start=4
startretries=3
autorestart=true
redirect_stderr=true
loglevel=debug


[program:websoc7]
command=python3 server.py -p 9007
directory=/home/ullget/dist/ullget/newsocket/server_socket/
user=root
stdout_logfile=/home/ullget/dist/ullget/newsocket/server_socket/logs/server.log
numprocs=1
numprocs_start=4
startretries=3
autorestart=true
redirect_stderr=true
loglevel=debug

[program:websoc8]
command=python3 server.py -p 9008
directory=/home/ullget/dist/ullget/newsocket/server_socket/
user=root
stdout_logfile=/home/ullget/dist/ullget/newsocket/server_socket/logs/server.log
numprocs=1
numprocs_start=4
startretries=3
autorestart=true
redirect_stderr=true
loglevel=debug

